# etc/opendias.service
[Unit]
Description=openDias, document imaging archive service
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/opendias
PIDFile=/var/run/opendias.pid
KillSignal=USR1

[Install]
WantedBy=multi-user.target
